# icmp-redirect/victim/_system/etc/rc.local
#!/bin/bash
route add -net 192.168.2.0/24  gw 192.168.1.10

#accept redirect
echo 1 > /proc/sys/net/ipv4/conf/all/accept_redirects
